# agent-maintain.service
[Unit]
Description=Agent System Cognitive Maintenance
Documentation=file:///home/frost/.gemini/memory/maintain.py

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/frost/.gemini/memory/maintain.py
StandardOutput=append:/home/frost/.cortexdb/maintain.log
StandardError=append:/home/frost/.cortexdb/maintain.log

# Don't kill agent processes if maintenance takes long
TimeoutStartSec=120
