# agent-maintain.service
[Unit]
Description=Agent Memory Maintenance (periodic)

[Service]
Type=oneshot
ExecStart=%h/.gemini/memory/maintain.py
StandardOutput=journal
StandardError=journal
