# agent-maintain.timer
[Unit]
Description=Agent System Maintenance Timer (every 6 hours)
Documentation=file:///home/frost/.gemini/memory/maintain.py

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
# Fire immediately on wake if missed during sleep/hibernate
Persistent=true

[Install]
WantedBy=timers.target
