# agent-maintain.timer
[Unit]
Description=Agent Memory Maintenance (periodic trigger)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
